# s6/service/pipewire-pulse/run
#!/bin/execlineb -P

fdmove -c 2 1
if { s6-svwait -U ../dbus-session-bus ../pipewire }
if { s6-svwait -u ../pipewire-media-session }
emptyenv -c
s6-envdir ../dbus-session-bus/env pipewire-pulse
