# scripts-system/etc/cron.daily/fstrim
#!/bin/sh
#
# Script to trim SSDs that support it

/sbin/fstrim --all --quiet
