# configs-user/.local/share/s6/service/.s6-svscan/finish
#!/bin/execlineb -P

echo User's supervision subtree has been successfully brought down!
