# s6/service/pulseaudio/run
#!/bin/execlineb -P

fdmove -c 2 1
if { s6-svwait -U ../dbus-session-bus }
s6-envdir ../dbus-session-bus/env
emptyenv -c pulseaudio --daemonize=no --high-priority=yes --realtime=yes --log-target=stderr --use-pid-file=no
